# backup.timer
[Unit]
Description=Running backup service every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
